# deploy/systemd/lain-mckenna.service
[Unit]
Description=Lain Terence McKenna Character Server
After=network.target lain-wired.service
PartOf=lain.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/local-lain
ExecStartPre=/bin/bash -c 'fuser -k 3004/tcp 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'mkdir -p /root/.lain-mckenna/workspace && [ -d workspace/characters/mckenna ] && cp -r workspace/characters/mckenna/* /root/.lain-mckenna/workspace/ || true'
ExecStart=/usr/bin/node dist/index.js mckenna --port 3004
Restart=on-failure
RestartSec=5

Environment=NODE_ENV=production
Environment=LAIN_HOME=/root/.lain-mckenna
Environment=CHARACTER_PROVIDER=google
Environment=CHARACTER_MODEL=gemini-2.5-flash
Environment=CHARACTER_API_KEY_ENV=GOOGLE_API_KEY
Environment=WIRED_LAIN_URL=http://localhost:3000
Environment=PORT=3004
EnvironmentFile=/opt/local-lain/.env
EnvironmentFile=/opt/local-lain/deploy/env/lain-mckenna.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lain-mckenna

[Install]
WantedBy=lain.target
