# deploy/systemd/lain.target
[Unit]
Description=Town - All Services
Wants=lain-gateway.service lain-telegram.service lain-voice.service

[Install]
WantedBy=multi-user.target
