# deploy/systemd/lain-pkd.service
[Unit]
Description=Lain Philip K. Dick Character Server
After=network.target lain-wired.service
PartOf=lain.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/local-lain
ExecStartPre=/bin/bash -c 'fuser -k 3003/tcp 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'mkdir -p /root/.lain-pkd/workspace && [ -d workspace/characters/pkd ] && cp -r workspace/characters/pkd/* /root/.lain-pkd/workspace/ || true'
ExecStart=/usr/bin/node dist/index.js pkd --port 3003
Restart=on-failure
RestartSec=5

Environment=NODE_ENV=production
Environment=LAIN_HOME=/root/.lain-pkd
Environment=CHARACTER_PROVIDER=anthropic
Environment=CHARACTER_MODEL=claude-haiku-4-5-20251001
Environment=CHARACTER_API_KEY_ENV=ANTHROPIC_API_KEY
Environment=WIRED_LAIN_URL=http://localhost:3000
Environment=PORT=3003
EnvironmentFile=/opt/local-lain/.env
EnvironmentFile=/opt/local-lain/deploy/env/lain-pkd.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lain-pkd

[Install]
WantedBy=lain.target
