# deploy/systemd/lain-dr-claude.service
[Unit]
Description=Lain Dr. Claude Character Server
After=network.target lain-web.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wired-lain
ExecStart=/usr/bin/node dist/index.js dr-claude --port 3002
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/wired-lain/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lain-dr-claude

[Install]
WantedBy=lain.target
